# packaging/debian/etc/notify-me/sv/notify-me/run
#!/bin/bash
exec 2>&1
exec chpst -u notifyme /opt/notify-me/run.sh

