# packaging/debian/etc/notify-me/sv/notify-me/run
#!/bin/bash

exec chpst -unotifyme /opt/notify-me/run.sh

